# meta/meta-mediatek-aud/recipes-core/busybox/busybox/busybox-klogd.service.in
[Unit]
Description=Kernel Logging Service

[Service]
ExecStart=/bin/klogd-start -n

[Install]
WantedBy=multi-user.target
