# meta/meta-mediatek-aud/recipes-core/busybox/busybox/tcpdump.service.in
[Unit]
Description=Tcpdump Logging Service

[Service]
ExecStart=/bin/tcpdump-start

[Install]
WantedBy=multi-user.target

