# meta/meta-mediatek-aud/recipes-apps/upgrade_app/upgrade_app.service
[Unit]
Description=upgrade_app

[Service]
ExecStart=/bin/upgrade_app
Type=simple

[Install]
Alias=upgrade
WantedBy=multi-user.target