# meta/meta-mediatek/recipes-core/udev/udev-only-hdmi/hdmi_hotplug.rules
ACTION=="change", KERNEL=="card0", SUBSYSTEM=="drm", RUN+="/etc/udev/scripts/hdmi_hotplug.sh"
